# suricata/log-shipper.service
[Unit]
Description=Suricata Log Shipper - Envía alertas al Dashboard Central
After=network.target suricata.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/html/pruebasoftwarelibre/suricata
ExecStart=/usr/bin/python3 /var/www/html/pruebasoftwarelibre/suricata/log_shipper.py
Restart=always
RestartSec=10

# Logs del servicio
StandardOutput=journal
StandardError=journal
SyslogIdentifier=log-shipper

[Install]
WantedBy=multi-user.target
